# prometheus/files/checkrestart_exporter.service
[Unit]
Description=Run the checkrestart prometheus exporter script

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/checkrestart_exporter.sh

